# tag-linux-laptop/local/share/systemd/user/bgrotate.service
[Unit]
Description=Desktop background rotation service

[Service]
Type=oneshot
ExecStart=%h/.local/bin/bgrotate
