# tag-linux-laptop/local/share/systemd/user/bgrotate.timer
[Unit]
Description=Desktop background rotation timer

[Timer]
OnCalendar=Mon *-*-* 00:00:00
Persistent=true

[Install]
WantedBy=timers.target

